# scripts/supervisor.lmis-dashboard.conf
[program:lmisdb]
command=node server/app.js
directory=/home/travisci/lmis-dashboard/latest/
user=travisci
autostart=true
autorestart=true
startretries=3
stdout_logfile=/home/travisci/lmis-dashboard/lmis-dashboard.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10
stderr_logfile=/home/travisci/lmis-dashboard/lmis-dashboard.log
stderr_logfile_maxbytes=1MB
stderr_logfile_backups=10
stopsignal=TERM
environment=NODE_ENV="production",IP="127.0.0.1",COUCH_HOST="http://dev.lomis.ehealth.org.ng",SESSION_SECRET="replace this"
